# templates/upstart-salt-minion.conf
description "salt-minion"

start on (net-device-up
          and local-filesystems
          and runlevel [2345])
stop on runlevel [!2345]

respawn limit 10 5

exec  ${buildout:directory}/bin/salt-minion -c ${buildout:directory}
